# ansible/files/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/Maksim/db_config
User=Maksim
WorkingDirectory=/home/Maksim/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always

[Install]
WantedBy=multi-user.target
